# files/systemd/system/nix-daemon.service.d/override.conf
[Unit]
After=nix.mount
Requires=nix.mount